# etc/systemd/system/bluetooth-comm.path
[Unit]
Description=Determina si existe un dispositivo conectado en rfcomm0

[Path]
PathExists=/dev/rfcomm0

[Install]
WantedBy=multi-user.target
