# etc/systemd/system/stop-upload.service
[Unit]
Description=Stops data upload
Conflicts=upload-server.service

[Service]
Type=oneshot
ExecStart=/bin/upload-stop.sh
