# etc/systemd/system/upload-server.service
[Unit]
Description=Servicio asociado a la unidad upload-server.timer

[Service]
User=pi
ExecStart=upload-files.sh
