# etc/systemd/system/upload-server.timer
[Unit]
Description=Timer unit of upload-server

[Timer]
OnCalendar=*-*-* 17:31

[Install]
WantedBy=timers.target
